# airflow.service
[Unit]
Description=Apache Airflow Standalone
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=max
Group=max
WorkingDirectory=/home/max/stack

# .env is loaded first; Environment= lines below act as defaults / overrides
EnvironmentFile=-/home/max/stack/.env
Environment=AIRFLOW_HOME=/home/max/stack
Environment=PATH=/home/max/stack/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=AIRFLOW__API__PORT=8091
Environment=AIRFLOW__API__HOST=0.0.0.0
Environment=AIRFLOW__CORE__EXECUTION_API_SERVER_URL=http://localhost:8091/execution/
Environment=AIRFLOW__CORE__DAGS_FOLDER=/home/max/stack/dags
Environment=AIRFLOW__CORE__LOAD_EXAMPLES=False
Environment=AIRFLOW__CORE__MP_START_METHOD=fork

ExecStart=/home/max/stack/.venv/bin/airflow standalone
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=airflow

[Install]
WantedBy=multi-user.target
